# config/systemd/user/pulseaudio-monitor.service
[Unit]
Description=PulseAudio Monitor Service

[Service]
ExecStart=/home/ludwig/bin/monitor-pulseaudio

[Install]
WantedBy=default.target
